# utils/crontab
*  *  *  *  * /usr/local/bin/python3 /app/update.py
